# deploy/quranbot_admin_supervisor.conf
[program:quranbot-admin]
command=/home/www/code/quranbot-admin/.venv/bin/uvicorn main:app --app-dir src --workers=2 --port=8010
stdout_logfile=/home/www/code/quranbot-admin/logs/supervisor.log
directory=/home/www/code/quranbot-admin
user=www
autorestart=true
redirect_stderr=true
